# etc/periodic_generation.crontab
SHELL=/bin/bash
STAT_CONTAINER_DIR=/home/bernhara/PC_CONTAINERs/pc_to_cs_usage_hr_stat_generator
#
# CS BASIC
#
# every monday, stats of previous week
0 6 * * 1	 cd ${STAT_CONTAINER_DIR} && docker compose up cs2pc-basic-previous-week-stat-generator

# every month, stats of previous month
0 6 1 * *	 cd ${STAT_CONTAINER_DIR} && docker compose up cs2pc-basic-previous-month-stat-generator
#
# CS PLUS
#
# every monday, stats of previous week
0 6 * * 1	 cd ${STAT_CONTAINER_DIR} && docker compose up cs2pc-plus-previous-week-stat-generator

# every month, stats of previous month
0 6 1 * *	 cd ${STAT_CONTAINER_DIR} && docker compose up cs2pc-plus-previous-month-stat-generator
